# recipes-core/bryton860-core/files/bryton860-core.service
[Unit]
Description=Start bryton860-core application
After=weston.socket

[Service]
Environment="XDG_RUNTIME_DIR=/run/"
Environment="WAYLAND_DISPLAY=wayland-0"

ExecStart=/usr/bin/bryton860-core

[Install]
WantedBy=multi-user.target
